# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor.log

[program:nginx]
command=nginx
autostart=true
autorestart=true

[program:php-fpm]
command=php-fpm8.3 -F
autostart=true
autorestart=true
